# openrc/k8snet
#!/sbin/openrc-run
description="Sets up networking in firecracker"
iface="${iface:-eth0}"

depend() {
    provide net
    before hostname
}

start() {
    source /etc/cmdline

    if [ -n "$KOPT_hostname" ]; then
        ebegin "Writing /proc/cmdline hostname"
        echo "$KOPT_hostname" > /etc/hostname
        eend $?
    fi

    if [ -n "$KOPT_k8s_ip" ]; then
        ebegin "Setting $iface IP address from /proc/cmdline"
        ip link set dev "$iface" up && ip addr add "$KOPT_k8s_ip" dev "$iface"
        eend $?
    fi

    if [ -n "$KOPT_k8s_gw" ]; then
        ebegin "Setting default gateway from /proc/cmdline"
        ip route add default via "$KOPT_k8s_gw"
        eend $?
    fi

    if [ -n "$KOPT_resolvconf" ]; then
        ebegin "Writing /etc/resolv.conf from /proc/cmdline "
        echo "$KOPT_resolvconf" | base64 -d > /etc/resolv.conf
        eend $?
    fi
}
