# openrc/livenessd
#!/sbin/openrc-run

command="/usr/local/bin/livenessd"
command_args="${LIVENESSD_OPTIONS}"
command_background="true"
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
    need net
    after lxd-init
}
