# openrc/lxd-data
#!/sbin/openrc-run
description="Sets up LXD data volume"
dev="${dev:-/dev/vdb}"

depend() {
    before overlay lxd
}

do_mount() {
    ebegin "Mounting LXD data volume"
    mkdir -p /var/lib/lxd && mount -t ext4 "$dev" /var/lib/lxd
    eend $?
}
start() {
    if ! do_mount; then
        ebegin "Formatting $dev as ext4"
        mkfs.ext4 -L lxd "$dev"
        eend $?

        do_mount
    fi

    return 0
}
