# openrc/lxd-data
#!/sbin/openrc-run
description="Sets up LXD data volume"
device="${device:-/dev/vdb}"

depend() {
	before lxd
}

do_mount() {
    ebegin "Mounting LXD data volume"
    mkdir -p /var/lib/lxd && mount -t ext4 "$device" /var/lib/lxd
    eend $?
}
start() {
    if ! do_mount; then
        ebegin "Formatting $device as ext4"
        mkfs.ext4 -L lxd "$device"
        eend $?

        do_mount
    fi
}
