# openrc/noop
#!/sbin/openrc-run

start() {
    return 0
}
